# startup/blinky.service
[Unit]
Description=Luxafor Blinky Service for Ticket Dash
# Ensures network is up and ticket_watcher has started
After=network.target ticket-watcher.service
# Tries to start ticket-watcher.service if not already running
Wants=ticket-watcher.service

[Service]
User=integotec
WorkingDirectory=/home/integotec/ticket-dash

# Command to start the blinky.py script
# Ensure 'python3' is the correct interpreter and can find necessary libraries (e.g., pyluxafor)
# If using a virtualenv: /home/integotec/ticket-dash/pyenv/bin/python3
ExecStart=/usr/bin/python3 /home/integotec/ticket-dash/blinky.py

# Restart policy
Restart=always
RestartSec=10

# Standard output and error logging
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
# The script reads from ./tickets.
ReadOnlyPaths=/home/integotec/ticket-dash/tickets

[Install]
WantedBy=multi-user.target
